# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:redis]
command=redis-server --dir /data --dbfilename dump.rdb
autostart=true
autorestart=true
stderr_logfile=/var/log/redis.err.log
stdout_logfile=/var/log/redis.out.log

[program:backend]
# We assume the app object is in backend.main:app
command=uvicorn backend.main:app --host 0.0.0.0 --port 8000 --reload
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/backend.err.log
stdout_logfile=/var/log/backend.out.log
environment=PYTHONUNBUFFERED=1

[program:celery]
# We assume the celery app is in backend.celery_worker:celery_app
# Using pool=solo for easier debugging in some envs, but prefork is better for production.
# Given single container and heavy GPU tasks, we might want concurrency=1 to avoid OOM on GPU.
command=celery -A backend.celery_worker.celery_app worker --loglevel=info --concurrency=1
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/celery.err.log
stdout_logfile=/var/log/celery.out.log
environment=PYTHONUNBUFFERED=1
